# etc/systemd/monitor.service
[Unit]
Description=Real-Time Server Monitoring Agent
After=network.target

[Service]
ExecStart=/usr/local/bin/monitor-realtime.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target

